# crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.dtdOzc/crontab installed on Mon Feb 16 14:04:00 2026)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Picframe cron jobs

# Fotos synchronisieren: Alle 30 Min zwischen 06:00 und 22:59 Uhr
*/30 6-22 * * * /home/schmali/picframe/scripts/sync_photos.sh >> /home/schmali/picframe/scripts/sync.log 2>&1

# WLAN Nachtschaltung
0 23 * * * sudo rfkill block wifi
0 6 * * * sudo rfkill unblock wifi

